# compose-cd.timer
[Unit]
Description=Continuous Deployment for docker-compose

[Timer]
OnUnitActiveSec=1min
Unit=compose-cd.service

[Install]
WantedBy=timers.target
