# compose-cd-cleanup.timer
[Unit]
Description=compose-cd cleanup timer

[Timer]
OnCalendar=daily
Unit=compose-cd-cleanup.service

[Install]
WantedBy=default.target
