# compose-cd-cleanup.service
[Unit]
Description=compose-cd cleanup

[Service]
Type=oneshot
ExecStart=/usr/bin/compose-cd cleanup

[Install]
WantedBy=default.target
